# auto-deployment/roles/diamond/templates/spvd.conf.j2

[program:diamond]
command=/usr/local/bin/diamond -c {{working_directory}}/diamond/diamond.conf
autorestart=true
directory={{working_directory}}/diamond